# resources/supervisor/benthos.conf
[program:benthos]
command = $GOPATH/bin/benthos
autostart = true
stdout_logfile = /var/log/benthos-out.log
stderr_logfile = /var/log/benthos-err.log
